# nantpclient.timer
[Unit]
Description=A non accurate NTP client
After=network.target

[Timer]
OnBootSec=15sec
OnUnitActiveSec=1d

[Install]
WantedBy=timers.target
