# nantpclient.timer
[Unit]
Description=A non accurate NTP client
After=network.target

[Timer]
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
